# desktop/crontab
@reboot ~/irclog-autogit.sh
@reboot ~/backup.sh
